# deploy/tonysgolfy-backend.service
[Unit]
Description=tonysgolfy backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/tonysgolfy/backend
Environment=APP_HOST=127.0.0.1
Environment=APP_PORT=3000
ExecStart=/srv/tonysgolfy/backend/target/release/backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
